# assets/ccs-proxy.service
[Unit]
Description=CCSwitch Proxy Server
After=network.target
Documentation=https://github.com/pomeluce/ccswitch

[Service]
ExecStart=/usr/bin/ccs proxy serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
